# data/xdg-desktop-portal-gnome.service.in
[Unit]
Description=Portal service (GNOME implementation)
After=graphical-session.target
Requisite=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.gnome
ExecStart=@libexecdir@/xdg-desktop-portal-gnome
